# crontab
# Weekly Report Generator - Cron Schedule
# Runs every Friday at 12:00 PM (noon) Eastern Time
# 
# Format: minute hour day-of-month month day-of-week command
# Day of week: 0=Sunday, 5=Friday

0 12 * * 5 cd /app && python main.py --mode auto >> /var/log/cron.log 2>&1

# Empty line required at end of crontab file
